# data/rauc.service.in
[Unit]
Description=Rauc Update Service
Documentation=https://rauc.readthedocs.io

[Service]
Type=dbus
BusName=de.pengutronix.rauc
ExecStart=@bindir@/rauc --mount=/run/rauc service
RuntimeDirectory=rauc
MountFlags=slave
