# data/rauc.service.in
[Unit]
Description=Rauc Update Service

[Service]
Type=dbus
BusName=de.pengutronix.rauc
ExecStart=@bindir@/rauc service
