# data/rauc.service
[Unit]
Description=Rauc Update Service

[Service]
Type=dbus
BusName=de.pengutronix.rauc
ExecStart=@BINDIR@/rauc service
